# etc/systemd/system/video-wall.service
# systemd service file - For Linux service management
# /etc/systemd/system/video-wall.service

[Unit]
Description=Video Wall Control System
After=network.target

[Service]
Type=simple
User=videowall
WorkingDirectory=/opt/video-wall
Environment=PATH=/opt/video-wall/venv/bin
ExecStart=/opt/video-wall/venv/bin/python app.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
